# misc/crontab
RAINVIEEW_APPROOT=$HOME/Rainview
PLACK_ENV=real

2 * * * * cd $RAINVIEW_APPROOT && ./script/env carton exec ./script/import_precip.pl 

